# deploy/systemd/cerious-gateway.service
[Unit]
Description=Cerious Gateway C++ Service
After=network-online.target cerious-exchange.service
Wants=network-online.target
Requires=cerious-exchange.service

[Service]
Type=simple
WorkingDirectory=/opt/cerious-systems
EnvironmentFile=/etc/cerious/cerious.env
ExecStart=/opt/cerious-systems/native/gateway-cpp/build/cerious_gateway --host ${CERIOUS_BACKEND_HOST} --port ${CERIOUS_BACKEND_PORT} --execution-host ${CERIOUS_EXCHANGE_HOST} --execution-port ${CERIOUS_EXCHANGE_HTTP_PORT} --root /opt/cerious-systems
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
